# prism/scheduling/com.prism.external-feed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prism.external-feed</string>
    <key>ProgramArguments</key>
    <array>
        <string>$HOME/work/prism/.venv/bin/prism</string>
        <string>process-external-feeds</string>
    </array>
    <key>WorkingDirectory</key>
    <string>$HOME/work/prism</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>$HOME/work/prism/data/launchd-external-feed.log</string>
    <key>StandardErrorPath</key>
    <string>$HOME/work/prism/data/launchd-external-feed.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>$HOME/work/prism/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
